# loginvsid
#!/bin/bash
### BEGIN INIT INFO
# Provides:          loginvsid
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: loginvsid
# Description:       loginvsid
### END INIT INFO

start() {
  cd /loginvsi
  env $(cat .env | grep ^[A-Z] | xargs) docker stack deploy --compose-file ./docker-compose.yml VSI
}

stop() {
  echo "Nothing to stop, yet"   
}

case $1 in
  start|stop) "$1" ;;
esac