# loginvsid
#!/bin/bash
### BEGIN INIT INFO
# Provides:          loginvsid
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: loginvsid
# Description:       loginvsid
### END INIT INFO

start() {
  cd /loginvsi
  docker stack deploy --compose-file ./docker-compose.yml VSI
}

stop() {
  #cd /loginvsi
  #docker stack deploy --compose-file ./docker-compose.yml VSI
}

case $1 in
  start|stop) "$1" ;;
esac